# infrastructure/ansible/roles/data-service/templates/data-service.service.j2
[Unit]
Description=AdminLTE Data Service
After=network.target postgresql.service redis.service adminlte-auth-service.service
Requires=postgresql.service redis.service

[Service]
Type=simple
User={{ adminlte_user | default('adminlte') }}
Group={{ adminlte_user | default('adminlte') }}
EnvironmentFile=/opt/adminlte/data-service/.env
WorkingDirectory=/opt/adminlte/data-service
ExecStart=/usr/bin/java -jar data-service.jar
ExecReload=/bin/kill -s HUP $MAINPID
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=adminlte-data-service

[Install]
WantedBy=multi-user.target